# nbavision.service
# NBAVision webhook server — systemd unit (Linux)
# Install: sudo cp nbavision.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable nbavision
#          sudo systemctl start nbavision
#
# Edit paths and user below.

[Unit]
Description=NBAVision Webhook Server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/NBAVisionEngine
Environment="PATH=/path/to/NBAVisionEngine/.venv/bin:/usr/bin"
Environment="NBAVISION_SECRET=your_strong_secret_here"
Environment="NBAVISION_BROWSER_PROFILE=/path/to/NBAVisionEngine/browser_profile"
ExecStart=/path/to/NBAVisionEngine/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
